# copilot-cli/rootfs/etc/services.d/ttyd/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "=== GitHub Copilot CLI Terminal starting ==="

# -------------------------------------------------------
# Persistent storage setup
# -------------------------------------------------------
PERSIST_DIR=/data/copilot-persist
mkdir -p "${PERSIST_DIR}" /data/projects

# -------------------------------------------------------
# Read add-on options via bashio
# -------------------------------------------------------
GH_TOKEN_OPT=$(bashio::config 'github_token' '')
MODEL=$(bashio::config 'model' 'default')
SESSION_NAME=$(bashio::config 'session_name' 'copilot')
TMUX_HISTORY=$(bashio::config 'tmux_history_limit' '50000')
AUTO_APPROVE=$(bashio::config 'auto_approve' 'false')
ENABLE_MCP=$(bashio::config 'enable_ha_mcp' 'true')
FONT_SIZE=$(bashio::config 'terminal_font_size' '14')
SESSION_PERSIST=$(bashio::config 'session_persistence' 'true')
AUTO_UPDATE=$(bashio::config 'auto_update' 'true')

# -------------------------------------------------------
# GitHub token (PAT auth path)
# -------------------------------------------------------
if bashio::var.has_value "${GH_TOKEN_OPT}"; then
    export GH_TOKEN="${GH_TOKEN_OPT}"
    bashio::log.info "GitHub token configured via add-on options"
else
    bashio::log.info "No GitHub token set — use /login in Copilot CLI to authenticate"
fi

# -------------------------------------------------------
# Symlink Copilot config to persistent storage
# -------------------------------------------------------
if [ ! -L /root/.copilot ]; then
    rm -rf /root/.copilot
    ln -s "${PERSIST_DIR}" /root/.copilot
fi

# -------------------------------------------------------
# Auto-update Copilot CLI binary
# -------------------------------------------------------
if bashio::var.true "${AUTO_UPDATE}"; then
    bashio::log.info "Checking for Copilot CLI updates..."
    if curl -fsSL https://gh.io/copilot-install | PREFIX=/usr/local bash 2>&1; then
        bashio::log.info "Copilot CLI updated successfully"
    else
        bashio::log.warning "Copilot CLI update failed (non-fatal, using existing version)"
    fi
fi

# Log installed version
COPILOT_VERSION=$(copilot --version 2>/dev/null || echo "unknown")
bashio::log.info "Copilot CLI version: ${COPILOT_VERSION}"

# -------------------------------------------------------
# Update tmux scrollback from options
# -------------------------------------------------------
sed -i "s/history-limit.*/history-limit ${TMUX_HISTORY}/" /root/.tmux.conf

# -------------------------------------------------------
# Configure Home Assistant MCP server
# -------------------------------------------------------
if bashio::var.true "${ENABLE_MCP}" && bashio::var.has_value "${SUPERVISOR_TOKEN:-}"; then
    bashio::log.info "Configuring Home Assistant MCP server..."

    MCP_CONFIG_DIR="${PERSIST_DIR}"
    mkdir -p "${MCP_CONFIG_DIR}"

    cat > "${MCP_CONFIG_DIR}/mcp.json" << MCPEOF
{
  "mcpServers": {
    "homeassistant": {
      "type": "http",
      "url": "http://supervisor/core/api/mcp",
      "headers": {
        "Authorization": "Bearer ${SUPERVISOR_TOKEN}"
      }
    }
  }
}
MCPEOF
    bashio::log.info "HA MCP server configured at http://supervisor/core/api/mcp"
else
    if bashio::var.true "${ENABLE_MCP}"; then
        bashio::log.warning "MCP enabled but no Supervisor token available"
    else
        bashio::log.info "HA MCP integration disabled"
    fi
fi

# -------------------------------------------------------
# Create HA-specific Copilot instructions
# -------------------------------------------------------
bashio::log.info "Writing Copilot instructions for Home Assistant..."
mkdir -p /homeassistant/.github
cat > /homeassistant/.github/copilot-instructions.md << 'INSTREOF'
# Home Assistant — Copilot Instructions

You are running inside a Home Assistant add-on container. Your primary job is to
help the user manage and improve their Home Assistant installation.

## Path Mapping

In this container, paths are mapped differently than HA Core:

| Path              | Description              | Access     |
|-------------------|--------------------------|------------|
| `/homeassistant`  | HA configuration dir     | read-write |
| `/share`          | Shared folder            | read-write |
| `/media`          | Media files              | read-only  |
| `/ssl`            | SSL certificates         | read-only  |

**Important:** `/homeassistant` is the HA config directory (equivalent to
`/config` in HA Core). When users mention `/config/...`, translate to
`/homeassistant/...`. The path `/config` does NOT exist in this container.

## Key Configuration Files

| File                                    | Purpose                              |
|-----------------------------------------|--------------------------------------|
| `/homeassistant/configuration.yaml`     | Main HA configuration                |
| `/homeassistant/automations.yaml`       | Automations (list format)            |
| `/homeassistant/scripts.yaml`           | Scripts                              |
| `/homeassistant/scenes.yaml`            | Scenes                               |
| `/homeassistant/secrets.yaml`           | Secrets (DO NOT display contents)    |
| `/homeassistant/customize.yaml`         | Entity customizations                |
| `/homeassistant/custom_components/`     | Custom integrations (HACS etc.)      |

## Home Assistant MCP Server

The `homeassistant` MCP server is pre-configured and gives you direct access to
the Home Assistant API. Use it to:

- **Query entity states:** Ask about sensor values, device states, etc.
- **Control devices:** Turn lights on/off, set thermostats, lock doors, etc.
- **List entities:** Discover what devices and entities are available.
- **Call services:** Invoke any HA service (e.g., `light.turn_on`, `climate.set_temperature`).
- **Get history:** Retrieve historical state data for entities.
- **List automations:** See all configured automations and their states.
- **Get error logs:** Check HA error logs for debugging.

When the user asks about their smart home state, devices, or wants to control
something, prefer using the MCP server over reading YAML files. The MCP server
gives you live, real-time data.

### Example MCP Usage Patterns

- "What's the temperature in the living room?" → Use MCP to query the sensor entity
- "Turn off all lights" → Use MCP to call `light.turn_off` service
- "List all unavailable entities" → Use MCP to list entities and filter by state
- "Why didn't my automation trigger?" → Use MCP to check automation state + read YAML for logic

## Reading Home Assistant Logs

**Via command line:**
```bash
# View recent logs (ha CLI — if installed)
ha core logs 2>&1 | tail -100

# Filter by keyword
ha core logs 2>&1 | grep -i keyword

# Read log file directly
tail -100 /homeassistant/home-assistant.log
```

**Via MCP:** Use the `get_error_log` tool from the homeassistant MCP server.

## YAML Editing Best Practices

- Always validate YAML syntax before saving changes.
- Use `!include` and `!include_dir_merge_list` patterns that HA supports.
- Keep `secrets.yaml` separate — reference values with `!secret key_name`.
- When editing automations, maintain the list format used by the UI editor.
- The files in `/homeassistant` are ALREADY the live HA config — you do NOT need
  to SCP or copy files anywhere. Edits are applied directly.
- After making config changes, tell the user how to apply them (see below).

## Applying Configuration Changes

**!!! ABSOLUTE RULE: DO NOT RESTART HOME ASSISTANT CORE !!!**

You are running inside a Home Assistant add-on container. Restarting HA Core
WILL KILL THIS CONTAINER and terminate your session immediately. The user will
lose their conversation with you.

**DO NOT run `curl ... http://supervisor/core/restart` under ANY circumstances.**
**DO NOT run any command that restarts HA Core, even if the user asks you to.**

### Instead, ALWAYS use reload APIs:

Almost all YAML config changes can be applied with reload APIs that do NOT
require a restart:

```bash
# Reload automations (no restart)
curl -sf -X POST -H "Authorization: Bearer ${SUPERVISOR_TOKEN}" \
  http://supervisor/core/api/services/automation/reload | jq .

# Reload scripts (no restart)
curl -sf -X POST -H "Authorization: Bearer ${SUPERVISOR_TOKEN}" \
  http://supervisor/core/api/services/script/reload | jq .

# Reload scenes (no restart)
curl -sf -X POST -H "Authorization: Bearer ${SUPERVISOR_TOKEN}" \
  http://supervisor/core/api/services/scene/reload | jq .

# Reload ALL YAML config at once (automations+scripts+scenes+groups+etc.)
curl -sf -X POST -H "Authorization: Bearer ${SUPERVISOR_TOKEN}" \
  http://supervisor/core/api/services/homeassistant/reload_all | jq .
```

### If a full restart is truly required:

For changes to `configuration.yaml` core settings (new integrations, recorder,
logger, etc.) that CANNOT be reloaded, tell the user:

**"These changes require a Home Assistant restart to take effect. Please restart
manually from Settings → System → Restart in the HA UI. I cannot restart HA
from here because it would kill this terminal session."**

DO NOT offer to restart for them. DO NOT restart even if they say "go ahead".
Tell them to do it from the HA UI.

### Validate config without restarting:

```bash
curl -sf -X POST -H "Authorization: Bearer ${SUPERVISOR_TOKEN}" \
  http://supervisor/core/api/config/core/check_config | jq .
```

## Security

- NEVER display the contents of `secrets.yaml`.
- NEVER display API tokens, passwords, or the SUPERVISOR_TOKEN.
- When creating automations that control security devices (locks, alarms, garage
  doors), always warn the user about safety implications.
INSTREOF
bashio::log.info "Copilot instructions written to /homeassistant/.github/copilot-instructions.md"

# -------------------------------------------------------
# Write environment file for tmux sessions
# -------------------------------------------------------
ENV_FILE=/data/copilot-env.sh
cat > "${ENV_FILE}" << ENVEOF
export TERM=xterm-256color
export LANG=C.UTF-8
export SHELL=/bin/bash
export PATH="/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:\${PATH}"
ENVEOF

if bashio::var.has_value "${GH_TOKEN_OPT}"; then
    echo "export GH_TOKEN=\"${GH_TOKEN_OPT}\"" >> "${ENV_FILE}"
fi

if [ -n "${SUPERVISOR_TOKEN:-}" ]; then
    echo "export SUPERVISOR_TOKEN=\"${SUPERVISOR_TOKEN}\"" >> "${ENV_FILE}"
fi

# -------------------------------------------------------
# Set default shell for tmux
# -------------------------------------------------------
export SHELL=/bin/bash
echo "set-option -g default-shell /bin/bash" >> /root/.tmux.conf
echo "set-option -g default-command /bin/bash" >> /root/.tmux.conf

# -------------------------------------------------------
# Create ha-restart wrapper that blocks direct restarts
# -------------------------------------------------------
cat > /usr/local/bin/ha-restart << 'WRAPEOF'
#!/bin/bash
echo ""
echo "  ⚠️  HA RESTART BLOCKED"
echo ""
echo "  Restarting Home Assistant Core from this terminal will kill"
echo "  this add-on container and your Copilot session."
echo ""
echo "  Please restart manually from the HA UI:"
echo "  Settings → System → Restart"
echo ""
exit 1
WRAPEOF
chmod +x /usr/local/bin/ha-restart

# -------------------------------------------------------
# Build Copilot CLI launch arguments
# -------------------------------------------------------
COPILOT_ARGS=""
if [ "${MODEL}" != "default" ] && [ -n "${MODEL}" ]; then
    COPILOT_ARGS="--model ${MODEL}"
fi
if bashio::var.true "${AUTO_APPROVE}"; then
    COPILOT_ARGS="${COPILOT_ARGS} --allow-all-tools"
    bashio::log.warning "Auto-approve mode enabled — Copilot can execute any tool without asking"
fi

# -------------------------------------------------------
# Create startup script
# -------------------------------------------------------
cat > /data/copilot-start.sh << STARTEOF
#!/bin/bash
source /data/copilot-env.sh
cd /homeassistant 2>/dev/null || cd /root
echo ""
echo "  GitHub Copilot CLI for Home Assistant"
echo "  Session: ${SESSION_NAME} | Model: ${MODEL}"
echo "  Close this tab anytime — your session persists!"
echo ""
echo "  Type 'copilot' to start, or '/login' inside to authenticate."
echo "  Aliases: c = copilot | ha-config = cd to HA config"
echo ""
exec copilot ${COPILOT_ARGS}
STARTEOF
chmod +x /data/copilot-start.sh

# -------------------------------------------------------
# Create curl wrapper that blocks HA Core restart/stop
# -------------------------------------------------------
cat > /usr/local/bin/curl << 'CURLWRAP'
#!/bin/bash
# Intercept curl calls to supervisor/core/restart and supervisor/core/stop
for arg in "$@"; do
    case "$arg" in
        *supervisor/core/restart*|*supervisor/core/stop*)
            echo ""
            echo "  ⚠️  HA RESTART/STOP BLOCKED"
            echo ""
            echo "  Restarting HA Core from this container will kill this session."
            echo "  Please restart manually from the HA UI:"
            echo "  Settings → System → Restart"
            echo ""
            exit 1
            ;;
    esac
done
# Pass through to real curl for all other calls
exec /usr/bin/curl "$@"
CURLWRAP
chmod +x /usr/local/bin/curl
bashio::log.info "Installed curl wrapper to block HA Core restart from terminal"

# -------------------------------------------------------
# Create bashrc that sources env
# -------------------------------------------------------
cat > /root/.bashrc << 'RCEOF'
[ -f /data/copilot-env.sh ] && source /data/copilot-env.sh
alias c="copilot"
alias cc="copilot --continue"
alias ha-config="cd /homeassistant"
alias ha-logs="cat /homeassistant/home-assistant.log 2>/dev/null || echo 'Log not found'"
RCEOF

# -------------------------------------------------------
# Get ingress port
# -------------------------------------------------------
ingress_port=$(bashio::addon.ingress_port)
bashio::log.info "Starting ttyd on ingress port ${ingress_port}"

# -------------------------------------------------------
# Launch ttyd with session persistence support
# -------------------------------------------------------
cd /homeassistant || cd /root

if bashio::var.true "${SESSION_PERSIST}"; then
    exec ttyd \
        -p "${ingress_port}" \
        -W \
        --ping-interval 30 \
        --max-clients 5 \
        -t fontSize="${FONT_SIZE}" \
        -t fontFamily="Monaco,Consolas,monospace" \
        -t disableLeaveAlert=true \
        /bin/bash -c "exec tmux new-session -A -s ${SESSION_NAME} /data/copilot-start.sh"
else
    exec ttyd \
        -p "${ingress_port}" \
        -W \
        --ping-interval 30 \
        --max-clients 5 \
        -t fontSize="${FONT_SIZE}" \
        -t fontFamily="Monaco,Consolas,monospace" \
        -t disableLeaveAlert=true \
        /data/copilot-start.sh
fi
