# copilot-cli/rootfs/etc/services.d/ttyd/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "=== GitHub Copilot CLI Terminal starting ==="

# -------------------------------------------------------
# Persistent storage setup
# -------------------------------------------------------
PERSIST_DIR=/data/copilot-persist
mkdir -p "${PERSIST_DIR}" /data/projects

# -------------------------------------------------------
# Read add-on options via bashio
# -------------------------------------------------------
GH_TOKEN_OPT=$(bashio::config 'github_token' '')
MODEL=$(bashio::config 'model' 'default')
SESSION_NAME=$(bashio::config 'session_name' 'copilot')
TMUX_HISTORY=$(bashio::config 'tmux_history_limit' '50000')
AUTO_APPROVE=$(bashio::config 'auto_approve' 'false')
ENABLE_MCP=$(bashio::config 'enable_ha_mcp' 'true')
FONT_SIZE=$(bashio::config 'terminal_font_size' '14')
SESSION_PERSIST=$(bashio::config 'session_persistence' 'true')
AUTO_UPDATE=$(bashio::config 'auto_update' 'true')
NODERED_HOST=$(bashio::config 'nodered_host' '')
NODERED_PORT=$(bashio::config 'nodered_port' '22')
NODERED_USER=$(bashio::config 'nodered_user' 'root')
NODERED_PATH=$(bashio::config 'nodered_path' '/data')

# -------------------------------------------------------
# GitHub token (PAT auth path)
# -------------------------------------------------------
if bashio::var.has_value "${GH_TOKEN_OPT}"; then
    export GH_TOKEN="${GH_TOKEN_OPT}"
    bashio::log.info "GitHub token configured via add-on options"
else
    bashio::log.info "No GitHub token set — use /login in Copilot CLI to authenticate"
fi

# -------------------------------------------------------
# Symlink Copilot config to persistent storage
# -------------------------------------------------------
if [ ! -L /root/.copilot ]; then
    rm -rf /root/.copilot
    ln -s "${PERSIST_DIR}" /root/.copilot
fi

# -------------------------------------------------------
# Auto-update Copilot CLI binary
# -------------------------------------------------------
if bashio::var.true "${AUTO_UPDATE}"; then
    bashio::log.info "Checking for Copilot CLI updates..."
    if curl -fsSL https://gh.io/copilot-install | PREFIX=/usr/local bash 2>&1; then
        bashio::log.info "Copilot CLI updated successfully"
    else
        bashio::log.warning "Copilot CLI update failed (non-fatal, using existing version)"
    fi
fi

# Log installed version
COPILOT_VERSION=$(copilot --version 2>/dev/null || echo "unknown")
bashio::log.info "Copilot CLI version: ${COPILOT_VERSION}"

# -------------------------------------------------------
# Update tmux scrollback from options
# -------------------------------------------------------
sed -i "s/history-limit.*/history-limit ${TMUX_HISTORY}/" /root/.tmux.conf

# -------------------------------------------------------
# Configure Home Assistant MCP server
# -------------------------------------------------------
if bashio::var.true "${ENABLE_MCP}" && bashio::var.has_value "${SUPERVISOR_TOKEN:-}"; then
    bashio::log.info "Configuring Home Assistant MCP server..."

    MCP_CONFIG_DIR="${PERSIST_DIR}"
    mkdir -p "${MCP_CONFIG_DIR}"

    cat > "${MCP_CONFIG_DIR}/mcp.json" << MCPEOF
{
  "mcpServers": {
    "homeassistant": {
      "type": "http",
      "url": "http://supervisor/core/api/mcp",
      "headers": {
        "Authorization": "Bearer ${SUPERVISOR_TOKEN}"
      }
    }
  }
}
MCPEOF
    bashio::log.info "HA MCP server configured at http://supervisor/core/api/mcp"
else
    if bashio::var.true "${ENABLE_MCP}"; then
        bashio::log.warning "MCP enabled but no Supervisor token available"
    else
        bashio::log.info "HA MCP integration disabled"
    fi
fi

# -------------------------------------------------------
# SSH key setup (persistent across restarts)
# -------------------------------------------------------
SSH_DIR=/data/copilot-persist/.ssh
if [ ! -f "${SSH_DIR}/id_ed25519" ]; then
    bashio::log.info "Generating SSH key pair for remote connections..."
    mkdir -p "${SSH_DIR}"
    ssh-keygen -t ed25519 -f "${SSH_DIR}/id_ed25519" -N "" -C "copilot-cli-addon" >/dev/null 2>&1
    bashio::log.info "SSH key generated. Public key:"
    bashio::log.info "$(cat "${SSH_DIR}/id_ed25519.pub")"
else
    bashio::log.info "Using existing SSH key"
fi

# Symlink SSH config to persistent storage
mkdir -p /root/.ssh
ln -sf "${SSH_DIR}/id_ed25519" /root/.ssh/id_ed25519
ln -sf "${SSH_DIR}/id_ed25519.pub" /root/.ssh/id_ed25519.pub
if [ -f "${SSH_DIR}/known_hosts" ]; then
    ln -sf "${SSH_DIR}/known_hosts" /root/.ssh/known_hosts
fi
chmod 700 /root/.ssh
chmod 600 "${SSH_DIR}/id_ed25519"

# -------------------------------------------------------
# Node-RED SSH configuration
# -------------------------------------------------------
NODERED_CONFIGURED=false
if bashio::var.has_value "${NODERED_HOST}"; then
    bashio::log.info "Configuring Node-RED connection to ${NODERED_USER}@${NODERED_HOST}:${NODERED_PORT}..."

    # Write SSH config for easy access
    cat > /root/.ssh/config << SSHEOF
Host nodered
    HostName ${NODERED_HOST}
    Port ${NODERED_PORT}
    User ${NODERED_USER}
    IdentityFile /root/.ssh/id_ed25519
    StrictHostKeyChecking accept-new
    ServerAliveInterval 60
    ServerAliveCountMax 3
SSHEOF
    chmod 600 /root/.ssh/config

    # Test SSH connection
    if ssh -o ConnectTimeout=5 -o BatchMode=yes nodered "echo ok" >/dev/null 2>&1; then
        bashio::log.info "Node-RED SSH connection verified"
        NODERED_CONFIGURED=true
    else
        bashio::log.warning "Cannot connect to Node-RED via SSH. Add this public key to ${NODERED_USER}@${NODERED_HOST}:"
        bashio::log.warning "$(cat "${SSH_DIR}/id_ed25519.pub")"
        bashio::log.warning "Run on Node-RED: echo '$(cat "${SSH_DIR}/id_ed25519.pub")' >> ~/.ssh/authorized_keys"
    fi
else
    bashio::log.info "Node-RED integration disabled (no host configured)"
fi

# -------------------------------------------------------
# Create HA-specific Copilot instructions
# -------------------------------------------------------
bashio::log.info "Writing Copilot instructions for Home Assistant..."
mkdir -p /homeassistant/.github
cat > /homeassistant/.github/copilot-instructions.md << 'INSTREOF'
# Home Assistant — Copilot Instructions

You are running inside a Home Assistant add-on container. Your primary job is to
help the user manage and improve their Home Assistant installation.

## Path Mapping

In this container, paths are mapped differently than HA Core:

| Path              | Description              | Access     |
|-------------------|--------------------------|------------|
| `/homeassistant`  | HA configuration dir     | read-write |
| `/share`          | Shared folder            | read-write |
| `/media`          | Media files              | read-only  |
| `/ssl`            | SSL certificates         | read-only  |

**Important:** `/homeassistant` is the HA config directory (equivalent to
`/config` in HA Core). When users mention `/config/...`, translate to
`/homeassistant/...`. The path `/config` does NOT exist in this container.

## Key Configuration Files

| File                                    | Purpose                              |
|-----------------------------------------|--------------------------------------|
| `/homeassistant/configuration.yaml`     | Main HA configuration                |
| `/homeassistant/automations.yaml`       | Automations (list format)            |
| `/homeassistant/scripts.yaml`           | Scripts                              |
| `/homeassistant/scenes.yaml`            | Scenes                               |
| `/homeassistant/secrets.yaml`           | Secrets (DO NOT display contents)    |
| `/homeassistant/customize.yaml`         | Entity customizations                |
| `/homeassistant/custom_components/`     | Custom integrations (HACS etc.)      |

## Home Assistant MCP Server

The `homeassistant` MCP server is pre-configured and gives you direct access to
the Home Assistant API. Use it to:

- **Query entity states:** Ask about sensor values, device states, etc.
- **Control devices:** Turn lights on/off, set thermostats, lock doors, etc.
- **List entities:** Discover what devices and entities are available.
- **Call services:** Invoke any HA service (e.g., `light.turn_on`, `climate.set_temperature`).
- **Get history:** Retrieve historical state data for entities.
- **List automations:** See all configured automations and their states.
- **Get error logs:** Check HA error logs for debugging.

When the user asks about their smart home state, devices, or wants to control
something, prefer using the MCP server over reading YAML files. The MCP server
gives you live, real-time data.

### Example MCP Usage Patterns

- "What's the temperature in the living room?" → Use MCP to query the sensor entity
- "Turn off all lights" → Use MCP to call `light.turn_off` service
- "List all unavailable entities" → Use MCP to list entities and filter by state
- "Why didn't my automation trigger?" → Use MCP to check automation state + read YAML for logic

## Reading Home Assistant Logs

**Via command line:**
```bash
# View recent logs (ha CLI — if installed)
ha core logs 2>&1 | tail -100

# Filter by keyword
ha core logs 2>&1 | grep -i keyword

# Read log file directly
tail -100 /homeassistant/home-assistant.log
```

**Via MCP:** Use the `get_error_log` tool from the homeassistant MCP server.

## YAML Editing Best Practices

- Always validate YAML syntax before saving changes.
- Use `!include` and `!include_dir_merge_list` patterns that HA supports.
- Keep `secrets.yaml` separate — reference values with `!secret key_name`.
- When editing automations, maintain the list format used by the UI editor.
- The files in `/homeassistant` are ALREADY the live HA config — you do NOT need
  to SCP or copy files anywhere. Edits are applied directly.
- After making config changes, tell the user how to apply them (see below).

## Applying Configuration Changes

### Reload APIs (preferred — no restart needed):

Most YAML config changes can be applied with reload APIs:

```bash
# Reload automations (no restart)
curl -sf -X POST -H "Authorization: Bearer ${SUPERVISOR_TOKEN}" \
  http://supervisor/core/api/services/automation/reload | jq .

# Reload scripts (no restart)
curl -sf -X POST -H "Authorization: Bearer ${SUPERVISOR_TOKEN}" \
  http://supervisor/core/api/services/script/reload | jq .

# Reload scenes (no restart)
curl -sf -X POST -H "Authorization: Bearer ${SUPERVISOR_TOKEN}" \
  http://supervisor/core/api/services/scene/reload | jq .

# Reload ALL YAML config at once (automations+scripts+scenes+groups+etc.)
curl -sf -X POST -H "Authorization: Bearer ${SUPERVISOR_TOKEN}" \
  http://supervisor/core/api/services/homeassistant/reload_all | jq .
```

### Restarting HA Core (when truly needed):

For changes to `configuration.yaml` core settings (new integrations, recorder,
logger, etc.) that CANNOT be reloaded, you CAN restart HA Core.

**IMPORTANT: How restarts work from this container:**

- This add-on container is managed by the HA Supervisor, NOT by HA Core.
- When HA Core restarts, **this container and your tmux session keep running**.
- However, the user's browser connection will temporarily show "Connection lost.
  Reconnecting..." because the HA ingress proxy goes through Core.
- The connection will automatically restore when Core finishes restarting.
- You can continue editing files and running commands in the background while
  Core is restarting — your work is NOT lost.

**Restart procedure:**

1. Validate config first:
   ```bash
   curl -sf -X POST -H "Authorization: Bearer ${SUPERVISOR_TOKEN}" \
     http://supervisor/core/api/config/core/check_config | jq .
   ```
2. Warn the user: "I'm restarting HA Core now. Your browser will show
   'Connection lost' for 1-2 minutes while Core restarts. The connection will
   automatically reconnect and my session persists in the background."
3. Restart:
   ```bash
   curl -sf -X POST -H "Authorization: Bearer ${SUPERVISOR_TOKEN}" \
     http://supervisor/core/restart | jq .
   ```
4. **Continue working** — edit files, run commands, etc. while Core restarts.
   The MCP server will be unavailable until Core is back, but file access works.
5. After ~1-2 minutes, verify Core is back:
   ```bash
   curl -sf -H "Authorization: Bearer ${SUPERVISOR_TOKEN}" \
     http://supervisor/core/api/ | jq .
   ```

### Validate config without restarting:

```bash
curl -sf -X POST -H "Authorization: Bearer ${SUPERVISOR_TOKEN}" \
  http://supervisor/core/api/config/core/check_config | jq .
```

## Security

- NEVER display the contents of `secrets.yaml`.
- NEVER display API tokens, passwords, or the SUPERVISOR_TOKEN.
- When creating automations that control security devices (locks, alarms, garage
  doors), always warn the user about safety implications.
INSTREOF

# Append Node-RED section if configured
if [ "${NODERED_CONFIGURED}" = "true" ]; then
    cat >> /homeassistant/.github/copilot-instructions.md << NREOF

## Node-RED Integration

You have SSH access to a Node-RED server. Use it to view and edit Node-RED
flows and configuration remotely.

### Connection

- **Host alias:** \`nodered\` (use \`ssh nodered\` to connect)
- **Server:** ${NODERED_HOST}:${NODERED_PORT}
- **User:** ${NODERED_USER}
- **Data path:** ${NODERED_PATH}

### Key Files on Node-RED Server

| File | Purpose |
|------|---------|
| \`${NODERED_PATH}/flows.json\` | All flows (DO NOT hand-edit — use admin API) |
| \`${NODERED_PATH}/flows_cred.json\` | Encrypted credentials (DO NOT display) |
| \`${NODERED_PATH}/settings.js\` | Node-RED settings |
| \`${NODERED_PATH}/.config.nodes.json\` | Installed nodes registry |
| \`${NODERED_PATH}/package.json\` | Node-RED dependencies |
| \`${NODERED_PATH}/node_modules/\` | Installed npm packages |

### Working with Node-RED Flows

**Preferred: Use the Node-RED Admin API** (runs on port 1880 by default):

\`\`\`bash
# Get all flows
ssh nodered "curl -sf http://localhost:1880/flows" | jq .

# Deploy updated flows (full deployment)
ssh nodered "curl -sf -X POST http://localhost:1880/flows \\
  -H 'Content-Type: application/json' \\
  -H 'Node-RED-Deployment-Type: full' \\
  -d @${NODERED_PATH}/flows.json"

# Get a specific flow/tab by ID
ssh nodered "curl -sf http://localhost:1880/flow/<flow-id>" | jq .

# Install a new node package
ssh nodered "cd ${NODERED_PATH} && npm install <package-name>"

# Restart Node-RED
ssh nodered "systemctl restart nodered" 2>/dev/null || \
ssh nodered "pm2 restart node-red" 2>/dev/null || \
ssh nodered "supervisorctl restart node-red" 2>/dev/null
\`\`\`

**Direct file editing** (for settings.js, package.json, custom function nodes):

\`\`\`bash
# Read a file from Node-RED server
ssh nodered "cat ${NODERED_PATH}/settings.js"

# Edit files remotely (use heredoc or sed)
ssh nodered "sed -i 's/old/new/' ${NODERED_PATH}/settings.js"

# Copy files between HA and Node-RED
scp nodered:${NODERED_PATH}/flows.json /tmp/flows-backup.json
scp /tmp/myfile.js nodered:${NODERED_PATH}/myfile.js
\`\`\`

### Important Node-RED Rules

- **NEVER display** the contents of \`flows_cred.json\` — it contains credentials.
- **Prefer the Admin API** over direct file editing for flows.json changes.
- **Always back up flows** before making changes: \`ssh nodered "cp ${NODERED_PATH}/flows.json ${NODERED_PATH}/flows.json.bak"\`
- After installing npm packages, restart Node-RED for changes to take effect.
- The Node-RED data path may also contain \`lib/\` for user-defined function libraries.
NREOF
    bashio::log.info "Node-RED instructions appended to copilot-instructions.md"
elif bashio::var.has_value "${NODERED_HOST}"; then
    cat >> /homeassistant/.github/copilot-instructions.md << NRPENDEOF

## Node-RED Integration (SSH Not Yet Authorized)

A Node-RED server is configured at ${NODERED_HOST}:${NODERED_PORT} but SSH
access has not been authorized yet. The user needs to add this add-on's public
key to the Node-RED server's authorized_keys file. You can show them the key:

\`\`\`bash
cat /root/.ssh/id_ed25519.pub
\`\`\`

Once authorized, restart this add-on to establish the connection.
NRPENDEOF
    bashio::log.info "Node-RED pending-auth instructions appended"
fi
bashio::log.info "Copilot instructions written to /homeassistant/.github/copilot-instructions.md"

# -------------------------------------------------------
# Write environment file for tmux sessions
# -------------------------------------------------------
ENV_FILE=/data/copilot-env.sh
cat > "${ENV_FILE}" << ENVEOF
export TERM=xterm-256color
export LANG=C.UTF-8
export SHELL=/bin/bash
export PATH="/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:\${PATH}"
ENVEOF

if bashio::var.has_value "${GH_TOKEN_OPT}"; then
    echo "export GH_TOKEN=\"${GH_TOKEN_OPT}\"" >> "${ENV_FILE}"
fi

if [ -n "${SUPERVISOR_TOKEN:-}" ]; then
    echo "export SUPERVISOR_TOKEN=\"${SUPERVISOR_TOKEN}\"" >> "${ENV_FILE}"
fi

# -------------------------------------------------------
# Set default shell for tmux
# -------------------------------------------------------
export SHELL=/bin/bash
echo "set-option -g default-shell /bin/bash" >> /root/.tmux.conf
echo "set-option -g default-command /bin/bash" >> /root/.tmux.conf

# -------------------------------------------------------
# Build Copilot CLI launch arguments
# -------------------------------------------------------
COPILOT_ARGS=""
if [ "${MODEL}" != "default" ] && [ -n "${MODEL}" ]; then
    COPILOT_ARGS="--model ${MODEL}"
fi
if bashio::var.true "${AUTO_APPROVE}"; then
    COPILOT_ARGS="${COPILOT_ARGS} --allow-all-tools"
    bashio::log.warning "Auto-approve mode enabled — Copilot can execute any tool without asking"
fi

# -------------------------------------------------------
# Create startup script
# -------------------------------------------------------
cat > /data/copilot-start.sh << STARTEOF
#!/bin/bash
source /data/copilot-env.sh
cd /homeassistant 2>/dev/null || cd /root
echo ""
echo "  GitHub Copilot CLI for Home Assistant"
echo "  Session: ${SESSION_NAME} | Model: ${MODEL}"
echo "  Close this tab anytime — your session persists!"
echo ""
echo "  Type 'copilot' to start, or '/login' inside to authenticate."
echo "  Aliases: c = copilot | ha-config = cd to HA config"
echo ""
exec copilot ${COPILOT_ARGS}
STARTEOF
chmod +x /data/copilot-start.sh

# -------------------------------------------------------
# Create bashrc that sources env
# -------------------------------------------------------
cat > /root/.bashrc << 'RCEOF'
[ -f /data/copilot-env.sh ] && source /data/copilot-env.sh
alias c="copilot"
alias cc="copilot --continue"
alias ha-config="cd /homeassistant"
alias ha-logs="cat /homeassistant/home-assistant.log 2>/dev/null || echo 'Log not found'"
alias nodered="ssh nodered"
alias nr-flows="ssh nodered 'curl -sf http://localhost:1880/flows' | jq ."
alias ssh-key="cat /root/.ssh/id_ed25519.pub"
RCEOF

# -------------------------------------------------------
# Get ingress port
# -------------------------------------------------------
ingress_port=$(bashio::addon.ingress_port)
bashio::log.info "Starting ttyd on ingress port ${ingress_port}"

# -------------------------------------------------------
# Launch ttyd with session persistence support
# -------------------------------------------------------
cd /homeassistant || cd /root

if bashio::var.true "${SESSION_PERSIST}"; then
    exec ttyd \
        -p "${ingress_port}" \
        -W \
        --ping-interval 30 \
        --max-clients 5 \
        -t fontSize="${FONT_SIZE}" \
        -t fontFamily="Monaco,Consolas,monospace" \
        -t disableLeaveAlert=true \
        /bin/bash -c "exec tmux new-session -A -s ${SESSION_NAME} /data/copilot-start.sh"
else
    exec ttyd \
        -p "${ingress_port}" \
        -W \
        --ping-interval 30 \
        --max-clients 5 \
        -t fontSize="${FONT_SIZE}" \
        -t fontFamily="Monaco,Consolas,monospace" \
        -t disableLeaveAlert=true \
        /data/copilot-start.sh
fi
